# crontab
0 0 3 * * python3 /app/main.py >> /dev/stdout 2>&1
